# debian/daemon-runtime.cron.d
* * * * * root /etc/daemon-manager/watcher.sh lenta-indexer